# custom_look/custom_theme.service
[Unit]
Description=Theme Customizer

[Service]
#User=root
ExecStart=/usr/local/bin/custom_theme.sh

[Install]
WantedBy=graphical.target
