# custom_look/global_theme.timer
[Unit]
Description=Time for global theme customizer 

[Timer]
OnBootSec=15sec
OnUnitActiveSec=15sec

[Install]
WantedBy=timers.target