# custom_look/global_theme.service
[Unit]
Description=Global theme customizer on timer 

[Service]
Type=oneshot
ExecStart=/home/wissem/bin/global_theme
